# .root/etc/udev/rules.d/99-power-profiles.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/usr/bin/auto-power-profile"
