# .root/etc/systemd/system/power-check.service
[Unit]
Description=Check battery level and set power profile

[Service]
Type=oneshot
ExecStart=/usr/bin/auto-power-profile
