# .root/etc/systemd/system/power-check.timer
[Unit]
Description=Run power check every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
